# KEEPRUNNING/update-reboot.service
[Unit]
Description=Update and Reboot Service

[Service]
ExecStart=~/Documenten/KEEPRUNNING/update-reboot.sh

[Install]
WantedBy=default.target
